# examples/newspaper.service
[Unit]
Description=Feedmailer

[Service]
ExecStart=/usr/bin/newspaper-cronjob

[Install]
WantedBy=basic.target

